# homelab-vps/systemd/homelab-vps-apps.target
[Unit]
Description=Homelab VPS Applications
After=homelab-vps-infra.target
PartOf=homelab-vps.target

[Install]
WantedBy=homelab-vps.target
