# homelab-vps/systemd/homelab-vps-monitoring.target
[Unit]
Description=Homelab VPS Monitoring
After=homelab-vps-infra.target
PartOf=homelab-vps.target

[Install]
WantedBy=homelab-vps.target
