# homelab-vps/systemd/rkhunter.service
[Unit]
Description=Daily scheduled rkhunter update and scan
Documentation=https://wiki.ubuntuusers.de/rkhunter/
PartOf=homelab-vps-infra.target
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/homelab-vps/homelab-vps/scripts/rkhunter.sh
Environment=DISPATCHER_URL=http://100.69.1.1:5001
User=root

# Logging - integrates with Loki via config.alloy
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rkhunter

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
